# base/scripts/mysql/supervisor_mysql.conf
[program:mysqld]
command=pidproxy /u01/socket/mysqld.pid /init_run_mysqld.sh
process_name=%(program_name)s
directory=/tmp
autostart=true
autorestart=true
redirect_stderr=true
user=root
